# resources/pgadmin.service.in
[Unit]
Description=pgAdmin4 is web-based PostgreSQL management server.
Documentation=https://www.pgadmin.org/docs

[Service]
Type=exec
ExecStart=/bin/bash -l -c "@PGADMIN4_CMD@"

[Install]
WantedBy=default.target
